# api/docker/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info
logfile=/var/log/supervisor/system.log
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/php-fpm-out.log
stderr_logfile=/var/log/supervisor/php-fpm-err.log
user=www-data